# packages/twitch-embedder/twitch-embedder.service
[Unit]
Description=Twitch Chat Embedder Service for Semantic Memory
After=network.target

[Service]
Type=simple
User=codingbutter
WorkingDirectory=/home/codingbutter/GitHub/rusty-butter/packages/twitch-embedder
ExecStart=/usr/bin/node /home/codingbutter/GitHub/rusty-butter/packages/twitch-embedder/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables - Update these with your actual values
Environment="NODE_ENV=production"
Environment="OPENAI_API_KEY=${OPENAI_API_KEY}"
Environment="TWITCH_USERNAME=rustybutterbot"
Environment="TWITCH_OAUTH_TOKEN=oauth:haa94c2dvym9ioan8cdk32jrc23a9y"
Environment="TWITCH_CHANNEL=codingbutter"
Environment="CHAT_EMBEDDER_INTERVAL=30000"

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/codingbutter/GitHub/rusty-butter/packages/twitch-embedder

[Install]
WantedBy=multi-user.target